# roles/proxmox/files/pvebackup.timer
[Unit]
Description=pvebackup timer

[Timer]
OnCalendar=daily
RandomizedDelaySec=12h
Unit=pvebackup.service

[Install]
WantedBy=timers.target